# practicas/so/practica6/script-namespace.service
[Unit]
Description=A simple service
[Service]
ExecStart=/usr/local/bin/script.sh
InaccessibleDirectories=/home