# practicas/so/practica6/script.service
[Unit]
Description=A simple service
[Service]
ExecStart=/home/rul/bin/script.sh
